# .config/systemd/user/auto-update@.path
[Unit]
Description=Triggers the service that push org files on change
Documentation=man:systemd.path

[Path]
PathModified=%h/%I

[Install]
WantedBy=multi-user.target